# autostart/autotiling.desktop
[Desktop Entry]
Exec=/usr/bin/autotiling
Name=autotiling
OnlyStartIn=sway
Type=Application
