# autostart/gebaard.desktop
[Desktop Entry]
Exec=/usr/bin/gebaard
Name=gebaard
OnlyStartIn=sway
Categories=sway;i3
Type=Application
